# fritz-dyndns/bash-script-image/crontab
* 0-23 * * * sh ./home/cloudflare-ddns-updater.sh >> /logfile